# ansible/roles/grillgauge/templates/grillgauge-scan.service.j2
[Unit]
Description=GrillGauge BLE Device Scanner
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/grillgauge
Environment=PATH=/opt/grillgauge/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/grillgauge/venv/bin/python -m src.grillgauge.cli scan --timeout {{ grillgauge_scan_timeout_seconds }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grillgauge-scan
